# conf/etc/systemd/system/where-da-gas-at-front.service
[Unit]
Description=Vite React Application
After=network.target

[Service]
Type=simple
User=where-da-gas-at
WorkingDirectory=/var/www/where-da-gas-at/react-vite
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
